# sshd-totp-stop.service
[Unit]
Description=Stop sshd-totp

[Service]
ExecStart=/srv/bot/ssh-disable.sh
Type=simple

[Install]

